# scripts/rusty-hass-laptop.service
[Unit]
Description=[Service description]

[Service]
Restart=on-failure
RestartSec=5s
Type=simple
StandardOutput=journal
ExecStart=/home/rodriguess/.local/bin/rusty-hass-laptop

[Install]
WantedBy=default.target
